# sentricore-dns-proxy.service
[Unit]
Description=Sentricore DNS Proxy
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=sentricore
WorkingDirectory=/opt/sentricore-dns
Environment="PYTHONUNBUFFERED=1"
ExecStart=/opt/sentricore-dns/venv/bin/python app/dns/proxy.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
